# layout/System/Library/LaunchDaemons/com.insanj.ColendarWriter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

<dict>
	<key>Label</key>
	<string>com.insanj.ColendarWriter</string>

	<key>Program</key>
	<string>/usr/libexec/ColendarWriter</string>

	<key>UserName</key>
	<string>root</string>

	<key>RunAtLoad</key>
	<true/>
</dict>

</plist>
